# src/config_example/supervisord.conf
[unix_http_server]
file = supervisor.sock

[supervisord]
logfile = supervisord.log
# 10 MiB
logfile_maxbytes = 10485760
logfile_backups = 5
pidfile = supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix://supervisor.sock

[program:site_check_http]
command = site_check_http.py -c config_local/site_checker.yaml
numprocs = 1
autorestart = true
stdout_logfile = site_check_http.log
stderr_logfile = site_check_http_stderr.log


[program:site_check_db_writer.py]
command = site_check_db_writer.py -c config_local/site_checker.yaml
numprocs = 1
autorestart = true
stdout_logfile = site_check_db_writer.log
stderr_logfile = site_check_db_writer_stderr.log
